# com.excelkit.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.excelkit.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/lly/trae-project/ExcelKit/node_modules/.bin/webpack</string>
        <string>serve</string>
        <string>--mode</string>
        <string>development</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/lly/trae-project/ExcelKit</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/lly/Library/Logs/excelkit-server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/lly/Library/Logs/excelkit-server.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NODE_PATH</key>
        <string>/usr/local/lib/node_modules</string>
    </dict>
</dict>
</plist>
